# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_X11VNC_CREATE_GEOM)s -ac +extension GLX +render -noreset
environment=DISPLAY=:99
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

[program:xfce]
command=/usr/bin/startxfce4
environment=DISPLAY=:99,HOME=/root
autorestart=true
priority=200
startsecs=5
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce.log

[program:vnc]
command=/usr/bin/x11vnc -display :99 -nopw -forever -shared -rfbport 5900
autorestart=true
priority=300
startsecs=5
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc.log

[program:recording]
command=/usr/local/bin/systemd-monitor
environment=DISPLAY=:99,RECORDINGS_DIR=/recordings
autorestart=true
priority=400
startsecs=10
stdout_logfile=/var/log/supervisor/recording.log
stderr_logfile=/var/log/supervisor/recording.log
